# ddw.cron.example
#!/bin/bash

cd /path/to/deepin-daily-wallpaper

/path/to/pipenv run ./ddw.py 

